# deploy/supervisor-worker.conf
# /etc/supervisor/conf.d/sms-gateway-worker.conf
# Mantém o worker das filas sempre vivo (envio assíncrono de SMS).
# Ajusta o caminho /var/www/sms e o utilizador (ex.: www-data).

[program:sms-gateway-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/sms/artisan queue:work --queue=default --sleep=3 --tries=3 --backoff=30 --max-time=3600
directory=/var/www/sms
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/sms/storage/logs/worker.log
stopwaitsecs=3600
